# debian/postinst
#!/bin/sh
# postinst for igh-seeedstudio
set -e

case "$1" in
    configure)
        # --- INST-01: Install blacklist for stock Realtek drivers ---
        # Use "install /bin/true" pattern to prevent udev bypass
        cat > /etc/modprobe.d/blacklist-eth.conf <<'BLACKLIST'
# Prevent stock Realtek drivers from loading.
# EtherCAT uses ec_r8169 instead.
install r8169 /bin/true
install r8168 /bin/true
BLACKLIST

        # --- INST-03 + INST-04: Detect MAC and write ethercat.conf ---
        IFACE="enP8p1s0"
        SYSFS_PATH="/sys/class/net/${IFACE}/address"

        if [ -f "${SYSFS_PATH}" ]; then
            MAC=$(cat "${SYSFS_PATH}")
        else
            echo "WARNING: ${SYSFS_PATH} not found, cannot detect MAC address" >&2
            echo "ethercat.conf will need manual configuration" >&2
            MAC=""
        fi

        cat > /etc/ethercat.conf <<CONF
# EtherCAT Master configuration
# Generated by igh-seeedstudio postinst
MASTER0_DEVICE="${MAC}"
DEVICE_MODULES="r8169"
CONF

        # --- INST-02: Register kernel modules ---
        depmod -a || true

        # --- INST-05: Start service AFTER depmod ---
        if [ -d /run/systemd/system ]; then
            systemctl restart ethercat.service || true
        fi
        ;;
esac

# debhelper-generated code (systemctl enable, daemon-reload, triggers)
#DEBHELPER#

exit 0
